# dist/unix/systemd/qbittorrent-nox.service
[Unit]
Description=qBittorrent Daemon Service
After=network.target

[Service]
Type=forking
User=1000
ExecStart=qbittorrent-nox -d

[Install]
WantedBy=multi-user.target
